# tools/pritunl-zero-redirect.socket
[Unit]
Description=Pritunl Zero Redirect Server Socket

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=sockets.target
